# crontab
* * */10 * * /app/build/postmanpat reapmessages
