# crontab
* */12 * * * /app/build/postmanpat reapmessages
